# .config/systemd/user/raspberrypi-backup.timer
[Unit]
Description=Raspberry Pi Backup Timer

[Timer]
OnCalendar=*-*-01,15 12:00:00
Persistent=true

[Install]
WantedBy=timers.target
